# packer-course-bootstrap/myip.service
[Unit]
Description=myip service

[Service]
User=root
ExecStart=/golang-myip/build/server.x86_64
Restart=always

[Install]
WantedBy=multi-user.target
